# deploy/grab-sync/gp1-grab-sync.timer
[Unit]
Description=Run GP1 Grab Crane sync every five minutes

[Timer]
OnCalendar=*:0/5
Persistent=true
RandomizedDelaySec=15
Unit=gp1-grab-sync.service

[Install]
WantedBy=timers.target
